# schedules/v7.8/macmini/io.agentskills.auto-checkpoint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.agentskills.auto-checkpoint</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fengning/agent-skills/scripts/dx-job-wrapper.sh</string>
        <string>auto-checkpoint</string>
        <string>--</string>
        <string>/Users/fengning/agent-skills/scripts/auto-checkpoint.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/fengning/logs/dx/auto-checkpoint.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fengning/logs/dx/auto-checkpoint.launchd.log</string>
</dict>
</plist>
