# schedules/v7.8/macmini/io.agentskills.dx-janitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.agentskills.dx-janitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/agent-skills/scripts/dx-job-wrapper.sh</string>
        <string>dx-janitor</string>
        <string>--</string>
        <string>__HOME__/agent-skills/scripts/dx-janitor.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/logs/dx/dx-janitor.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/logs/dx/dx-janitor.launchd.log</string>
</dict>
</plist>
