# hive/node/hive-queen.service
[Unit]
Description=Hive Mind Queen Bee (Autonomous Dispatcher)
After=network.target

[Service]
Type=simple
# Assumes the repo is at ~/agent-skills. Use %h for home.
ExecStart=/usr/bin/python3 %h/agent-skills/hive/node/hive-queen.py
WorkingDirectory=%h/agent-skills
Restart=always
RestartSec=30
# Ensure we have the path to 'bd', 'cass', and standard tools
Environment="PATH=%h/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
