# systemd/dx-repo-memory-refresh.service
[Unit]
Description=DX Repo Memory Refresh
After=network-online.target

[Service]
Type=oneshot
Environment=PATH=%h/.local/share/mise/shims:%h/.local/share/mise/bin:%h/.local/bin:%h/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
ExecStart=%h/agent-skills/scripts/dx-repo-memory-refresh --repo-name agent-skills --repo-path %h/agent-skills --fallback-host macmini --push --create-pr
WorkingDirectory=%h
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.dx-state %h/.codex %h/.cache %h/.local %h/agent-skills %h/.config/gh %h/.ssh /tmp /private/tmp

[Install]
WantedBy=default.target
