# systemd/slack-coordinator.service
[Unit]
Description=Slack Coordination Service (V4.2.1 - Scoped Env)
After=network.target

[Service]
Type=simple
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin"

# V4.2.1 Security: Load Encrypted Credential with fallback
LoadCredentialEncrypted=op_token:%h/.config/systemd/user/op_token.cred
LoadCredential=op_token:%h/.config/systemd/user/op_token

# Use SCOPED env file (not shared ~/.agent-env)
# op run resolves op:// references at runtime
ExecStart=/bin/bash -c 'export OP_SERVICE_ACCOUNT_TOKEN="$(cat $CREDENTIALS_DIRECTORY/op_token)"; exec /home/linuxbrew/.linuxbrew/bin/op run --env-file=%h/.config/slack-coordinator/.env -- %h/agent-skills/slack-coordination/.venv/bin/python slack-coordinator.py'

Restart=on-failure
RestartSec=5
WorkingDirectory=%h/agent-skills/slack-coordination
MemoryMax=2G

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/slack-coordinator %h/.local/state/slack-coordination

[Install]
WantedBy=default.target
